# meteo.cron
0 * * * * /chemin de votre script/Extracteur_Météo.sh
